# autoemail/autoemail.conf
description "autoemail"
author "Goclis Yao <goclisyyh@gmail.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /path
exec /path/autoemail.py
respawn